# scripts/launchd/com.nicospencer.lenny.dogfood.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nicospencer.lenny.dogfood</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_ROOT_PLACEHOLDER/scripts/launchd/run-dogfood.sh</string>
  </array>

  <!-- Run every 6 hours (21600 seconds). -->
  <key>StartInterval</key>
  <integer>21600</integer>

  <!-- Persist across reboots. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>HOME_PLACEHOLDER/Library/Logs/lenny-dogfood/run.log</string>

  <key>StandardErrorPath</key>
  <string>HOME_PLACEHOLDER/Library/Logs/lenny-dogfood/run.log</string>

  <!-- Add Flutter SDK, Homebrew, and local binaries to PATH for dart, bd, fs. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>HOME_PLACEHOLDER/flutter/bin:HOME_PLACEHOLDER/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
